# deploy/launchd/ai.kakeya.grpc-runtime-prefill.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>ai.kakeya.grpc-runtime-prefill</string>
  <key>ProgramArguments</key><array>
    <string>/usr/bin/caffeinate</string><string>-dimsu</string>
    <string>/Users/fluffy314/.venv-distwan/bin/python</string>
    <string>/Users/fluffy314/Documents/Kakeya-LLM-Inference-engine-prefill-kv-network/scripts/start_grpc_runtime_server.py</string>
    <string>--backend</string><string>mlx</string>
    <string>--verifier-id</string><string>/Users/fluffy314/kakeya-models/gemma-4-26B-A4B-it-mlx-4bit</string>
    <string>--bind</string><string>127.0.0.1:51051</string>
    <string>--capacity</string><string>1</string>
    <string>--sink</string><string>4</string>
    <string>--window</string><string>2048</string>
    <string>--max-concurrent-rpcs</string><string>4</string>
    <string>--skip-cache-check</string>
    <string>--enable-prefill-cache</string>
    <string>--prefill-cache-gb</string><string>1</string>
    <string>--peer</string><string>169.254.27.104:53051</string>
    <string>--cache-peer</string><string>169.254.27.104:53051</string>
    <string>--cache-model-id</string><string>gemma-4-26B-A4B-it-mlx-4bit</string>
    <string>--model-revision</string><string>local-4bit-v1</string>
    <string>--tokenizer-revision</string><string>gemma4-v1</string>
    <string>--cache-quantization</string><string>4bit-mlx</string>
    <string>--cache-kv-dtype</string><string>bfloat16</string>
    <string>--cache-block-tokens</string><string>64</string>
    <string>--cache-tenant-id</string><string>private-fleet</string>
    <string>--fleet-psk-file</string><string>/Users/fluffy314/.kakeya/fleet.psk</string>
    <string>--node-id</string><string>head-runtime</string>
    <string>--advertise</string><string>169.254.187.239:51051</string>
    <string>--cache-advertise</string><string>169.254.187.239:51051</string>
    <string>--network-label</string><string>thunderbolt</string>
    <string>--network-priority</string><string>100</string>
    <string>--measured-rtt-ms</string><string>0.55</string>
    <string>--remote-prefill-min-tokens</string><string>128</string>
    <string>--cache-link-mbps</string><string>10000</string>
    <string>--cache-default-rtt-ms</string><string>0.55</string>
    <string>--prefill-min-savings-ratio</string><string>0</string>
    <string>--primary-prefill-penalty-ms</string><string>1000</string>
    <string>--network-http-host</string><string>127.0.0.1</string>
    <string>--network-http-port</string><string>8090</string>
    <string>--network-api-key</string><string>__NETWORK_KEY__</string>
    <string>--network-state</string><string>/Users/fluffy314/.kakeya/inference_network.json</string>
    <string>--network-telemetry-url</string><string>http://127.0.0.1:8090/v1/network/telemetry/tokens</string>
    <string>--network-telemetry-api-key</string><string>__NETWORK_KEY__</string>
    <string>--log-level</string><string>INFO</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/fluffy314/Documents/Kakeya-LLM-Inference-engine-prefill-kv-network</string>
  <key>EnvironmentVariables</key><dict>
    <key>PATH</key><string>/Users/fluffy314/.venv-distwan/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PYTHONPATH</key><string>/Users/fluffy314/Documents/Kakeya-LLM-Inference-engine-prefill-kv-network:/Users/fluffy314/Documents/Kakeya-LLM-Inference-engine-prefill-kv-network/sdks/python</string>
    <key>HF_HUB_OFFLINE</key><string>0</string>
    <key>TRANSFORMERS_OFFLINE</key><string>0</string>
  </dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ProcessType</key><string>Interactive</string>
  <key>StandardOutPath</key><string>/Users/fluffy314/.kakeya/grpc-runtime-prefill.log</string>
  <key>StandardErrorPath</key><string>/Users/fluffy314/.kakeya/grpc-runtime-prefill.log</string>
</dict></plist>
